# will.conf
description "W.I.L.L"
author  "Will Beddow <will@willbeddow.com>"

start on runlevel [234]
stop on runlevel [0156]

chdir /usr/local/will
exec /usr/local/will/main.py

respawn